# roles/certificates/templates/lego.timer
[Unit]
Description=Renewal of certificate from Let's Encrypt

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
RequiredBy=nginx.service
